# deploy/next.service
[Unit]
Description=Next.js Application Service
After=network.target

[Service]
User=your_username
WorkingDirectory=/path/to/nextjs/app
EnvironmentFile=/path/to/nextjs/.env
ExecStartPre=/usr/bin/npm install
ExecStart=/usr/bin/npm start
Restart=always

[Install]
WantedBy=multi-user.target
